# operation/mongodb.jp.cron
6 3 * * * root /usr/local/cockatoo/operation/rotate.sh
8 3 * * * root /usr/local/cockatoo/operation/loggz.sh

